# deploy/senal-mcp.service
[Unit]
Description=Daily Tracker MCP server (solo lectura, para Claude)
After=network.target

[Service]
Type=simple
User=senal
Group=senal
WorkingDirectory=/opt/senal
Environment=MCP_HTTP=1
Environment=MCP_HOST=127.0.0.1
Environment=MCP_PORT=8765
ExecStart=/opt/senal/.venv/bin/python /opt/senal/mcp_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
